# Voron/Macros/Sleep_Wake/klipper-watcher.service
[Unit]
Description=Klipper State Watcher
After=network.target

[Service]
Type=simple
User=pi                  # Use the user under which you run the script (e.g., 'pi')
ExecStart=/usr/bin/python3 /home/pi/klipper-watcher.py   # Path to your Python script
WorkingDirectory=/home/pi
Restart=on-failure
RestartSec=5             # Restart the service after 5 seconds in case of failure

[Install]
WantedBy=multi-user.target
